# systemd/rent-radar-map.service
[Unit]
Description=Rent Radar — Mapa interactivo (HTTP)
After=network.target

[Service]
User=YOUR_USER
WorkingDirectory=YOUR_PROJECT_DIR
EnvironmentFile=YOUR_PROJECT_DIR/.env
ExecStart=YOUR_PROJECT_DIR/venv/bin/python run_dashboard.py --serve --port 8080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
